# site/profile/files/jira/TELUS_jira_control
#!/bin/bash
# chkconfig: 345 99 99

# <Update with real code and delete this line>